# server/crontab
# UTC daily at 00:00
0 0 * * * /usr/bin/curl -fsS -H "x-cron-secret: ${CRON_SECRET}" http://127.0.0.1:4000/cron/refresh
